# service/walldisplay.service
[Unit]
Description=WallDisplay Sonos monitor and RGB matrix driver
Documentation=https://github.com/patojar/WallDisplay
Wants=network-online.target
After=network-online.target sound.target

[Service]
Type=simple
WorkingDirectory=/home/pi/WallDisplay
Environment=WALLDISPLAY_FLAGS=-display
EnvironmentFile=-/etc/default/wall-display
ExecStart=/home/pi/WallDisplay/bin/musicDisplay $WALLDISPLAY_FLAGS
Restart=on-failure
RestartSec=10s
StandardOutput=journal
StandardError=inherit

[Install]
WantedBy=multi-user.target
